# package_builder/bearouter-3.6.11-0.0.1/etc/rc.local
#!/bin/sh -e

# xl2tpd configurator
/usr/bin/cava.sh

exit 0
