# package_builder/candybox-3.6.11-0.0.1/etc/rc.local
#!/bin/sh -e

# candybox web configurator configurator
/usr/bin/candybox.sh

exit 0
